# deploy/praxis-backend.service
[Unit]
Description=Praxis backend (FastAPI async, uvicorn) - by CENAT
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/praxis/backend
EnvironmentFile=/opt/praxis/backend/.env
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.local/bin/uv run uvicorn app.main:app --host 127.0.0.1 --port 8040 --workers 1 --proxy-headers
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/praxis/backend/praxis-backend.log
StandardError=append:/opt/praxis/backend/praxis-backend.log

[Install]
WantedBy=multi-user.target
